# docker/supervisord.d/laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "if [ \"$SCHEDULER_ENABLED\" = \"true\" ] || [ \"$SCHEDULER_ENABLED\" = \"1\" ]; then while [ true ]; do (php /var/www/html/web-miniverse/artisan schedule:run --verbose --no-interaction &); sleep 60; done; else echo 'Scheduler is disabled. Set SCHEDULER_ENABLED=true to enable.'; sleep 3600; fi"
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
stderr_logfile=/var/log/supervisor/scheduler.err.log

