# docker/supervisord.d/laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "if [ \"$SCHEDULER_ENABLED\" = \"true\" ] || [ \"$SCHEDULER_ENABLED\" = \"1\" ]; then echo '[Scheduler] 已啟用 - 開始執行排程'; while true; do echo '[$(date +\"%%Y-%%m-%%d %%H:%%M:%%S\")] Running scheduler...'; php /var/www/html/web-miniverse/artisan schedule:run --verbose --no-interaction; sleep 60; done; else echo '[Scheduler] 未啟用 - 設置 SCHEDULER_ENABLED=true 以啟用'; sleep 3600; fi"
autostart=true
autorestart=true
numprocs=1
user=www-data
directory=/var/www/html/web-miniverse
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/scheduler.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

